# deploy/rfobserver.service
[Unit]
Description=RFObserver Sensor Pipeline
# The sensor is a NATS client to a remote RFS broker, so no local
# nats-server.service ordering is needed.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rfobserver
Group=rfobserver
# Config lives in a writable .env in the state dir (not an EnvironmentFile under
# read-only /etc). WorkingDirectory makes that .env the one both the app reads on
# startup and the UI writes on toggle/apply, so runtime changes persist across
# restarts. Do NOT add an EnvironmentFile: its env vars would outrank the .env
# and silently override UI-persisted changes.
WorkingDirectory=/var/lib/rfobserver
ExecStart=/opt/rfobserver/venv/bin/rfobserver run
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rfobserver

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/rfobserver /tmp/rfobserver

[Install]
WantedBy=multi-user.target
